# bin/service/focus_last.service
[Unit]
Description=run focus last

[Service]
Type=simple
User=al
Group=al
ExecStart=/home/al/bin/focus-last.py
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target

